# root/etc/services.d/sonarrannounced/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	cd /app/sonarrannounced
	s6-setuidgid abc python3.5 /app/sonarrannounced/bot.py
