# root/etc/services.d/sonarrannounced/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc python /app/sonarrannounced/bot.py \
	--config_file=/config/settings.cfg --data_dir=/config/data
